# config/blink.service
[Unit]
Description=blink — background screen recording daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=__EXEC_START__
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
